# deployment/supervisord.conf
[program:eva]
command=python manage.py runserver 0.0.0.0:8000
directory=/app

[program:celery]
command=celery -A eva worker
directory=/app

[program:redis]
command=redis-server
